# data/deploy/local/supervisor/mail_rq.conf
[program:mall_rq]

command= python  manage.py rqworker default low
directory=/vagrant_data/back
user=vagrant
autostart=true
autorestart=true
stderr_logfile=/vagrant_data/logs/rq_errors.log
stdout_logfile=/vagrant_data/logs/rq_access.log
environment = DJANGO_SETTINGS_MODULE="settings.local"