# data/deploy/local/supervisor/mall_api.conf
[program:mall_api]

command=sh /vagrant_data/deploy/local/gunicorn/start.sh
directory=/vagrant_data
user=vagrant
autostart=true
autorestart=true
stderr_logfile=/vagrant_data/logs/gu_errors.log
stdout_logfile=/vagrant_data/logs/gu_access.log