# data/deploy/pro/supervisor/mall_rq.conf
[program:mall_rq]

command=/home/django/venv/bin/python  manage.py rqworker default low
directory=/home/django/mall/data/back
user=django
autostart=true
autorestart=true
stderr_logfile=/home/django/mall/data/logs/rq_errors.log
stdout_logfile=/home/django/mall/data/logs/rq_access.log
environment = DJANGO_SETTINGS_MODULE="settings.production"